# examples/infracanvas-agent.service
[Unit]
Description=InfraCanvas Infrastructure Discovery Agent
Documentation=https://github.com/bytestrix/InfraCanvas
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/infracanvas/agent.env
ExecStart=/usr/local/bin/infracanvas start
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=infracanvas-agent

[Install]
WantedBy=multi-user.target
